# systemd/alarm-worker.service
[Unit]
Description=통합알람감시체계 Worker
After=network.target

[Service]
Type=simple
User=dietpi
WorkingDirectory=/home/dietpi/alarm
ExecStart=/usr/local/bin/node /home/dietpi/alarm/node_modules/.bin/tsx src/worker/index.ts
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
